# repo/packages/b/busybox/port/1.35.0/inittab
# Copyright (c) 2006-2024, RT-Thread Development Team
#
# SPDX-License-Identifier: Apache-2.0

# Run getty & login
# console::respawn:/sbin/getty 115200 /dev/console

# or without login, a direct access to shell
console::respawn:-/bin/ash